# build/supervisor/cert-monitor.ini
[program:cert-monitor]
command=/usr/local/qcloud-tools/cert-monitor
stopsignal=QUIT
autorestart=true
autostart=true