# build/supervisor/http-server.ini
[program:http-server]
command=/usr/local/qcloud-tools/http-server
stopsignal=QUIT
autorestart=true
autostart=true